# crontab
# START CRON JOB
SHELL=/bin/sh
PATH=/usr/local/bin:/usr/bin:/bin
15 15 * * * echo 'hello'
# END CRON JOB