# crontab
# START CRON JOB
* * * * * python3 ${XL_IDP_PATH_DOCKER}/scripts/send2telegram.py
# END CRON JOB